# scripts/deploy/dailyapp.service
[Unit]
Description=DailyApp Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=dailyapp
Group=dailyapp
WorkingDirectory=/opt/dailyapp
ExecStart=/usr/bin/java -Xms512m -Xmx1g -XX:+UseG1GC -jar /opt/dailyapp/*.jar
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dailyapp

# 环境变量
Environment=SPRING_PROFILES_ACTIVE=production
Environment=JAVA_OPTS=-Djava.security.egd=file:/dev/./urandom

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dailyapp /opt/backups /var/log/dailyapp

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 